# nginx_php7/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:php]
command=/usr/sbin/php-fpm7 -F

